# packages/nxenv-nxcli/nxenv_nxcli-5.0.0.dev0-py3-none-any.whl/nxcli/config/templates/systemd/nxenv-nxcli-nxenv-long-worker.service
[Unit]
Description="{{ nxcli_name }}-nxenv-short-worker %I"
PartOf={{ nxcli_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ nxcli_cmd }} worker --queue long
StandardOutput=file:{{ nxcli_dir }}/logs/worker.log
StandardError=file:{{ nxcli_dir }}/logs/worker.error.log
WorkingDirectory={{ nxcli_dir }}
