# packages/nxenv-nxcli/nxenv_nxcli-5.0.0.dev0.tar.gz/nxenv_nxcli-5.0.0.dev0/nxcli/config/templates/systemd/nxenv-nxcli-redis.target
[Unit]
After=network.target
Wants={{ nxcli_name }}-redis-cache.service {{ nxcli_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
