# packages/nxenv-nxcli/nxenv_nxcli-5.0.0.dev0.tar.gz/nxenv_nxcli-5.0.0.dev0/nxcli/config/templates/systemd/nxenv-nxcli-web.target
[Unit]
After=network.target
Wants={{ nxcli_name }}-nxenv-web.service {{ nxcli_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
