# packages/nxenv-nxcli/nxenv_nxcli-5.0.0.dev0.tar.gz/nxenv_nxcli-5.0.0.dev0/nxcli/config/templates/systemd/nxenv-nxcli.target
[Unit]
After=network.target
Requires={{ nxcli_name }}-web.target {{ nxcli_name }}-workers.target {{ nxcli_name }}-redis.target

[Install]
WantedBy=multi-user.target
