# packages/kitchenai/kitchenai-0.16.0.tar.gz/kitchenai-0.16.0/deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
exec python kitchenai/__main__.py run
